# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── 主應用 ──────────────────────────────────────────────────
[program:gunicorn]
command=gunicorn --config gunicorn.conf.py "mes_dashboard:create_app()"
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1
stopwaitsecs=30

; ── RQ Workers ──────────────────────────────────────────────
[program:worker-trace]
command=rq worker %(ENV_TRACE_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
stopwaitsecs=60

[program:worker-reject]
command=rq worker %(ENV_REJECT_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
stopwaitsecs=60

[program:worker-msd]
command=rq worker %(ENV_MSD_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
stopwaitsecs=60

[program:worker-production-history]
command=rq worker %(ENV_PRODUCTION_HISTORY_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
stopwaitsecs=60

[program:worker-yield-alert]
command=rq worker %(ENV_YIELD_ALERT_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
stopwaitsecs=60

[program:worker-warmup]
command=rq worker %(ENV_WARMUP_WORKER_QUEUE)s --url %(ENV_REDIS_URL)s
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
stopwaitsecs=120
