# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── 主應用 ──────────────────────────────────────────────────
[program:gunicorn]
command=gunicorn --config gunicorn.conf.py "mes_dashboard:create_app()"
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1
stopwaitsecs=30

; ── RQ Worker Pool Manager ──────────────────────────────────
; Replaces static per-queue worker programs with a single demand-driven
; manager.  Workers are spawned on-demand and self-terminate after
; POOL_IDLE_TIME seconds of inactivity.
;
; Key tuning env vars (set in docker-compose or .env):
;   POOL_GLOBAL_MAX        total worker cap (Oracle connection budget)
;   POOL_MAX_PER_QUEUE     default per-queue maximum
;   POOL_IDLE_TIME         seconds a worker waits idle before exiting
;   POOL_MAX_WARMUP_QUEUE_NAME=1   keep warmup capped at 1 worker
;
[program:worker-pool-manager]
command=python -m mes_dashboard.core.worker_pool_manager
directory=/app
environment=PYTHONPATH="/app/src"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
; stopwaitsecs must exceed POOL_IDLE_TIME so the manager can drain workers
; cleanly before supervisord force-kills it.
stopwaitsecs=120
stopsignal=TERM
