# deploy/mes-dashboard-downtime-worker.service
[Unit]
Description=MES Dashboard Downtime Query Worker (RQ, Conda Runtime)
Documentation=https://github.com/your-org/mes-dashboard
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mes-dashboard
EnvironmentFile=-/opt/mes-dashboard/.env
Environment="PYTHONPATH=/opt/mes-dashboard/src"

# RQ is pinned rq>=1.16.0,<2.0.0 which has NO `rq worker --job-execution-timeout`
# flag; per-job timeout is enforced at enqueue via job_timeout. Removed the invalid flag.
ExecStart=/usr/bin/env bash -lc 'exec "${CONDA_BIN:-/opt/miniconda3/bin/conda}" run --no-capture-output -n "${CONDA_ENV_NAME:-mes-dashboard}" rq worker "${DOWNTIME_WORKER_QUEUE:-downtime-query}" --url "${REDIS_URL:-redis://localhost:6379/0}"'

KillSignal=SIGTERM
TimeoutStopSec=60
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=mes-dashboard-downtime-worker

# Memory protection: downtime worker handles large Oracle analytical queries.
# MemoryMax prevents single large job from killing the VM.
MemoryHigh=3G
MemoryMax=4G

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/mes-dashboard/logs /opt/mes-dashboard/tmp

[Install]
WantedBy=multi-user.target
